# etc/services.d/confd/run
#!/usr/bin/with-contenv bash

/opt/confd --watch --backend env